# systemd/async-tcp-udp-server.service
# systemd/async-tcp-udp-server.service
[Unit]
Description=Async TCP/UDP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=async-server
Group=async-server

Environment=SERVER_PORT=8080

ExecStart=/usr/local/bin/async_tcp_udp_server
ExecStop=/bin/kill -TERM $MAINPID

# Более агрессивные таймауты
TimeoutStopSec=10
TimeoutStartSec=10
KillMode=process
KillSignal=SIGTERM

# Важно: НЕ перезапускать после остановки
Restart=no

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=async-tcp-udp-server

[Install]
WantedBy=multi-user.target