# systemd/async-tcp-udp-server.service
# Создай правильный unit файл
cat > systemd/async-tcp-udp-server.service << 'EOF'
[Unit]
Description=Async TCP/UDP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=async-server
Group=async-server

Environment=SERVER_PORT=8080
Environment=LOG_LEVEL=info

ExecStart=/usr/local/bin/async_tcp_udp_server

TimeoutStopSec=10
TimeoutStartSec=5
KillMode=mixed
KillSignal=SIGTERM

Restart=no

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=async-tcp-udp-server

[Install]
WantedBy=multi-user.target
EOF